# ops/pixelcraft-api.service
[Unit]
Description=PixelCraft Bloom API Server
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=pixelcraft
Group=pixelcraft
WorkingDirectory=/opt/pixelcraft-bloom/backend
Environment=PATH=/opt/pixelcraft-bloom/backend/.venv/bin
Environment=PYTHONPATH=/opt/pixelcraft-bloom/backend
EnvironmentFile=/opt/pixelcraft-bloom/backend/.env
ExecStart=/opt/pixelcraft-bloom/backend/.venv/bin/gunicorn app.main:app -w 4 -k uvicorn.workers.UvicornWorker -b 127.0.0.1:8000 --access-logfile /var/log/pixelcraft/access.log --error-logfile /var/log/pixelcraft/error.log
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/pixelcraft-bloom/backend /var/log/pixelcraft /tmp
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target